# cloudformation-openvidu/roles/supervisord/templates/openvidu.conf.j2
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=debug

[program:openvidu-server]
command=/usr/local/bin/openvidu-server.sh

{% if deploy_demos == "true" %}
[program:js-java]
{% if whichcert == "letsencrypt" or whichcert == "owncert" %}
command=/bin/bash /var/www/html/basic-webinar/app.sh {{ domain_name }}
{% else %}
command=/bin/bash /var/www/html/basic-webinar/app.sh {{ aws_public_hostname }}
{% endif %}
redirect_stderr=true

[program:classroom-demo]
{% if whichcert == "letsencrypt" or whichcert == "owncert" %}
command=/bin/bash /var/www/html/classroom/app.sh {{ domain_name }}
{% else %}
command=/bin/bash /var/www/html/classroom/app.sh {{ aws_public_hostname }}
{% endif %}
redirect_stderr=true
{% endif %}